# scripts/launchd/com.prav.wc26-preview.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.prav.wc26-preview</string>

  <!-- Wrapper script handles tournament-window + match-hour gating so this
       plist stays one-shot simple. -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/prav/Desktop/personal-projects/fifa-wc-26-prediction/scripts/launchd/run_if_tournament.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/prav/Desktop/personal-projects/fifa-wc-26-prediction</string>

  <!-- Fire every 15 minutes, 24/7. Wrapper exits in &lt;1s outside the
       tournament + match windows, so the wake-up cost is negligible. -->
  <key>StartInterval</key>
  <integer>900</integer>

  <!-- Surface PATH so vercel CLI (installed via Homebrew / npm-global) and
       python3 are findable when launchd runs us. -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/Users/prav/.npm-global/bin</string>
    <key>HOME</key>
    <string>/Users/prav</string>
  </dict>

  <!-- Capture launchd-side stdout/stderr (wrapper's own log is separate). -->
  <key>StandardOutPath</key>
  <string>/Users/prav/Desktop/personal-projects/fifa-wc-26-prediction/logs/launchd-stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/prav/Desktop/personal-projects/fifa-wc-26-prediction/logs/launchd-stderr.log</string>

  <!-- Don't relaunch on early-exit (wrapper's "skip" path returns 0). -->
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
